# scripts/frankenphp-panel.service
[Unit]
Description=FrankenPHP Panel
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=panel
Group=panel
WorkingDirectory=/opt/frankenphp-panel
EnvironmentFile=/opt/frankenphp-panel/.env
ExecStart=/opt/frankenphp-panel/frankenphp-panel
Restart=on-failure
RestartSec=5

# NoNewPrivileges=yes would block sudo (needed for site-create.sh); panel still runs as User=panel
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
